# dokuwiki/rootfs/etc/services.d/php-fpm/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: DokuWiki
# Runs the PHP-FPM daemon
# ==============================================================================
bashio::log.info "Starting PHP-FPM..."

export UPLOAD_LIMIT=64M

exec php-fpm83 --nodaemonize
